# srv_automation/imperative/ansible/roles/networking/networkmanager_migrate/templates/nm-migrate-rollback.timer.j2
# Managed by Ansible (networkmanager_migrate). Do not edit.
[Unit]
Description=Timer to rollback NetworkManager cutover if not confirmed

[Timer]
OnActiveSec={{ networkmanager_migrate_rollback_minutes }}min
AccuracySec=10s
Persistent=false

[Install]
WantedBy=timers.target

