# srv_automation/imperative/ansible/roles/streaming/kafka/configure/templates/kafka.service.j2
[Unit]
Description=Apache Kafka (KRaft)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User={{ kafka_user }}
Group={{ kafka_group }}
WorkingDirectory={{ kafka_install_dir }}

# Keep paths unambiguous: always use the symlinked install dir
ExecStart={{ kafka_install_dir }}/bin/kafka-server-start.sh {{ kafka_install_dir }}/config/kraft/server.properties
ExecStop={{ kafka_install_dir }}/bin/kafka-server-stop.sh

Restart=on-failure
RestartSec=5
LimitNOFILE=100000

# Logs go to journald unless you redirect inside Kafka config
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
